# shell/copilot-insights.service
[Unit]
Description=GitHub Copilot CLI Token Insights Dashboard Service
After=network.target

[Service]
Type=simple
WorkingDirectory=<PROJECT_DIR>
ExecStart=<PROJECT_DIR>/target/release/copilot-cli-insights
Restart=always
RestartSec=5
Environment=PORT=3000
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
